# vllm.service
[Unit]
Description=vLLM Server -- Ornith-1.0-35B FP8 on DGX Spark
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/vllm2
Environment=PORT=8000
Environment=HOST=0.0.0.0
Environment=MODEL_REPO=deepreinforce-ai/Ornith-1.0-35B-FP8
Environment=MODEL_PATH=/srv/vllm2/models/Ornith-1.0-35B-FP8
Environment=SERVED_MODEL_NAME=Ornith-1.0-35B
Environment=TOKENIZER=/srv/vllm2/models/Ornith-1.0-35B-FP8
Environment=MAX_MODEL_LEN=262144
Environment=MAX_NUM_BATCHED_TOKENS=8192
Environment=MAX_NUM_SEQS=3
Environment=GPU_MEM_UTIL=0.5
ExecStart=/srv/vllm2/vllm-server.sh
TimeoutStartSec=1800
TimeoutStopSec=120
KillSignal=SIGTERM
SendSIGKILL=yes
Restart=on-failure
RestartSec=15
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target
